# script/backup.service
[Unit]
Description=Run backup script

[Service]
ExecStart=/home/panda/projects/Melda/backup/backup.sh
